# deploy/systemd/pickletour-worker-general.service
[Unit]
Description=PickleTour Go General Worker
After=network.target mongod.service redis.service
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/pickletour
EnvironmentFile=/var/www/pickletour/.env
ExecStart=/var/www/pickletour/backend-go/bin/pickletour-worker-general
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
