# deploy/systemd/pickletour-scheduler.service
[Unit]
Description=PickleTour Go Scheduler
After=network.target mongod.service redis.service
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/pickletour
EnvironmentFile=/var/www/pickletour/.env
ExecStart=/var/www/pickletour/backend-go/bin/pickletour-scheduler
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
