# deploy/observer-vps/pickletour-observer.service
[Unit]
Description=PickleTour Observer Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=observer
Group=observer
WorkingDirectory=/var/lib/pickletour-observer
EnvironmentFile=/etc/pickletour/.env
ExecStart=/usr/local/bin/pickletour-observer
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
