# deploy/observer-vps/pickletour-observer.service
[Unit]
Description=PickleTour Observer Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/pickletour
EnvironmentFile=/etc/pickletour/observer.env
ExecStart=/usr/bin/env node backend/observerServer.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
